# app/.server/batch/schedule.cron
# Send Email
* * * * * cd /app && npm run batch common/SendEmail

# Cleanup unused S3 images
0 4 * * * cd /app && npm run batch common/RemoveUnusedS3Images