# app/.server/batch/schedule.cron
* * * * * cd /app && npm run batch common/SendEmail